# deploy/namvibe-daphne.service.example
[Unit]
Description=Namvibe Daphne ASGI server
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/namvibe
EnvironmentFile=/var/www/namvibe/.env
ExecStart=/var/www/namvibe/venv/bin/daphne -b 127.0.0.1 -p 8001 config.asgi:application
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
